# Collectors/TaskScheduler-HealthScore-V2.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Date>2025-01-01T00:00:00</Date>
    <Author>DBA Team</Author>
    <Description>Health Score V2 - Collectors Frecuentes (cada 5 minutos)</Description>
    <URI>\HealthScore\V2-Collectors-Frequent</URI>
  </RegistrationInfo>
  <Triggers>
    <CalendarTrigger>
      <Repetition>
        <Interval>PT5M</Interval>
        <StopAtDurationEnd>false</StopAtDurationEnd>
      </Repetition>
      <StartBoundary>2025-01-01T00:00:00</StartBoundary>
      <Enabled>true</Enabled>
      <ScheduleByDay>
        <DaysInterval>1</DaysInterval>
      </ScheduleByDay>
    </CalendarTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <UserId>S-1-5-21-XXXXXXXXXX-XXXXXXXXXX-XXXXXXXXXX-XXXX</UserId>
      <LogonType>Password</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>false</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>true</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <DisallowStartOnRemoteAppSession>false</DisallowStartOnRemoteAppSession>
    <UseUnifiedSchedulingEngine>true</UseUnifiedSchedulingEngine>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>PT10M</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>powershell.exe</Command>
      <Arguments>-ExecutionPolicy Bypass -File "C:\HealthScore\Collectors\Run-All.ps1" -Mode Frequent</Arguments>
      <WorkingDirectory>C:\HealthScore\Collectors</WorkingDirectory>
    </Exec>
  </Actions>
</Task>

<!-- 
INSTRUCCIONES DE USO:

1. Editar el archivo y cambiar:
   - UserId: Usar cuenta de servicio con permisos adecuados
   - Arguments: Ajustar la ruta al script Run-All.ps1
   - WorkingDirectory: Ajustar la ruta base

2. Crear 3 tareas:
   
   a) Frequent (cada 5 min):
      schtasks /Create /XML TaskScheduler-HealthScore-V2.xml /TN "HealthScore\V2-Collectors-Frequent"
   
   b) Periodic (cada 10 min):
      Copiar este XML, cambiar <Interval>PT10M</Interval> y Arguments a "-Mode Periodic"
      schtasks /Create /XML TaskScheduler-HealthScore-V2-Periodic.xml /TN "HealthScore\V2-Collectors-Periodic"
   
   c) Daily (cada 24h a las 2 AM):
      Copiar este XML, cambiar trigger a diario y Arguments a "-Mode Daily"
      schtasks /Create /XML TaskScheduler-HealthScore-V2-Daily.xml /TN "HealthScore\V2-Collectors-Daily"

3. Verificar:
   Get-ScheduledTask -TaskName "*HealthScore*"

4. Ejecutar manualmente para probar:
   Start-ScheduledTask -TaskName "HealthScore\V2-Collectors-Frequent"

NOTAS:
- La cuenta de servicio debe tener permisos de lectura en las instancias SQL monitoreadas
- La cuenta debe tener permisos de escritura en SQLNova.dbo.Inventario*Snapshot
- Recomendado: usar gMSA (Group Managed Service Account) en entornos de dominio
-->

